# yawu/Client/.service
[Unit]
Description=Package Manager Service
After=network.target

[Service]
Type=simple
ExecStart=/path/to/your/application/executable
Restart=on-failure
User=your_user

[Install]
WantedBy=multi-user.target
